# systemd/.config/systemd/user/old.xremap.service
[Unit]
Description=xremap

[Service]
; Restart=always
ExecStart=xremap /home/kyle/.config/xremap/config.yml --watch
